# daemontools.service
[Unit]
Description=daemontools
After=sysinit.target

[Service]
ExecStart=/usr/bin/svscanboot
Restart=always

[Install]
WantedBy=multi-user.target
