# daemontools.conf
# daemontools

start on runlevel [12345]
stop on runlevel [06]

respawn
exec /usr/bin/svscanboot
